# packages/andrzej-tools/boot-mail.service
[Unit]
Description=Send boot notification email
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment=CONFIG=/etc/andrzej-tools.conf
ExecStart=/usr/local/sbin/boot-mail

[Install]
WantedBy=multi-user.target
