# packages/andrzej-tools/ping-watch.service
[Unit]
Description=Ping watch (all targets from /etc/andrzej-tools.conf)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
LogLevelMax=3
Environment=CONFIG=/etc/andrzej-tools.conf
ExecStart=/usr/local/sbin/ping-watch
